# data/etc/network/if-up.d/iptables
#!/bin/bash
iptables-restore < /etc/iptables/iptables.rules