# data/etc/cron.weekly/disk-maintenance
#!/bin/bash
# search and destroy mac/windows garbage files
if [ -d /media ]
then
    find /media -iname "thumbs.db" -exec rm -rf {} \;
    find /media -iname ".ds_store" -exec rm -rf {} \;
    find /media -name '._*' -exec rm -rf {} \;
fi
# reduced io stress with weekly fstrim
for DEVICE in $(mount | grep ext4 | grep -v mapper | awk '{print $1}')
do
    fstrim "${DEVICE}"
done
# defragment ext4 devices
for DEVICE in $(mount | grep ext4 | awk '{print $1}')
do
    e4defrag "${DEVICE}"
done
