# data/etc/cron.weekly/system-updates
#!/bin/bash
# execute btrfs snapshot prior to running updates
# set headless flags for debian package updates

mv /etc/apt/sources.list /etc/apt/sources.list.bak
netselect-apt -sn -o /etc/apt/sources.list $(lsb_release -c | cut -f2)
aptitude clean
if ! aptitude update
then
    mv /etc/apt/sources.list.bak /etc/apt/sources.list
    aptitude clean
    aptitude update
fi
aptitude upgrade -yq
update-command-not-found