# data/etc/systemd/system/xboxdrv.service
[Unit]
Description=Xbox controller driver daemon

[Service]
Type=simple
User=root
ExecStart=/usr/bin/xboxdrv --detach-kernel-driver -c /etc/default/xboxdrv

[Install]
WantedBy=multi-user.target
