# data/template/etc/network/if-up.d/iptables
#!/bin/bash
iptables-restore < /etc/firewall.conf
