# data/template/etc/cron.weekly/aptitude
#!/bin/sh

# update packages weekly
aptitude clean
aptitude update
aptitude upgrade -y
update-command-not-found
