# data/template/etc/cron.weekly/e4defrag
#!/bin/sh

# defragment ext4 devices
for DEVICE in $(mount | grep ext4 | awk '{print $1}')
do
    e4defrag "${DEVICE}"
done
