# data/template/etc/cron.monthly/netselect
#!/bin/bash

# Update package mirrors
netselect-apt -s -n
aptitude clean
aptitude update
