# linux/debian/data/base/etc/cron.weekly/maintance
#!/bin/bash
/usr/local/bin/system-updates
/usr/local/bin/disk-maintenance

